# .config/systemd/user/battery-notify.service
# ~/.config/systemd/user/battery-notify.service
[Unit]
Description=Battery level notifier

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-notify
